# avsocks.service
[Unit]
Description=Fuck the gfw

[Service]
ExecStart=/usr/bin/avsocks

[Install]
WantedBy=multi-user.target
