# avsocks.socket
[Unit]
Description=Fuck the gfw

[Service]
ListenStream=127.0.0.1:4567

[Install]
WantedBy=sockets.target