# avsocks.socket
[Unit]
Description=Fuck the gfw

[Socket]
ListenStream=4567
Service=avsocks.service

[Install]
WantedBy=sockets.target