# newsocks.socket
[Unit]
Description=Fuck the gfw

[Service]
ListenStream=127.0.0.1:11080

[Install]
WantedBy=sockets.target